# visync.service
[Unit]
Description=visync
After=network.target
StartLimitIntervalSec=0
[Service]
Type=simple
Restart=always
RestartSec=1
User=david
ExecStart=/home/david/visync-dist/start-server

[Install]
WantedBy=multi-user.target
